# deploy/systemd/logrotate.timer
[Unit]
Description=Hourly rotation of log files

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
